# docker/rutorrent/rootfs/etc/services.d/nginx/run
#!/usr/bin/with-contenv bash
exec /usr/sbin/nginx -c /config/rutorrent-nginx/nginx.conf
