# docker/plexpy/rootfs/etc/services.d/plexpy/run
#!/usr/bin/with-contenv bash

exec s6-setuidgid swerp python /opt/plexpy/PlexPy.py --datadir /config
