# docker/deluge/rootfs/etc/services.d/deluged/run
#!/usr/bin/with-contenv bash

umask 002

exec s6-setuidgid swerp /usr/bin/deluged -c /config -d --loglevel=info -l /logs/deluge/deluged.log
