# docker/sickrage/rootfs/etc/services.d/sickrage/run
#!/usr/bin/with-contenv bash

umask 002

exec s6-setuidgid swerp python /app/sickrage/SickBeard.py --datadir /config
